# crontab
0 0 * * 2-6 /check_daily_job_line_calls.sh
5 0 * * 6 /email_weekly_job_line_report.sh
